# scripts/backup.service
[Unit]
Description="Multi Database Backup Service"
Documentation=https://github.com/semmidev/phylax
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=garuda
Group=garuda
ExecStart=/usr/local/bin/phylax -config /etc/phylax/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups /var/log/phylax /tmp

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
